# deployment/systemd/ai-airbagszentrum-deploy.service
[Unit]
Description=Deploy AiZentrum from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/deploy-ai-airbagszentrum
